# .ebextensions/cron/crontab
0 0 * * * /usr/bin/php /var/www/html/artisan paxifi:crontab
